# .config/autostart/compton.desktop
[Desktop Entry]
Comment=The window effects.
Exec=compton
Name=Compton
Type=Application
