# .config/autostart/tint2.desktop
[Desktop Entry]
Comment=The task bar.
Exec=tint2
Name=Tint2
Type=Application
